# launchagents/com.mlx.coding.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>

  <key>Label</key>
  <string>com.mlx.coding</string>

  <!-- Verify mlx_lm.server path after install: which mlx_lm.server -->
  <key>ProgramArguments</key>
  <array>
    <string>/usr/local/bin/mlx_lm.server</string>
    <string>--model</string>
    <string>mlx-community/Qwen2.5-Coder-14B-Instruct-4bit</string>
    <string>--host</string>
    <string>0.0.0.0</string>
    <string>--port</string>
    <string>8082</string>
  </array>

  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/var/log/mlx/coding.log</string>
  <key>StandardErrorPath</key>
  <string>/var/log/mlx/coding.error.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>

  <key>ThrottleInterval</key>
  <integer>10</integer>

</dict>
</plist>
