# launchagents/com.localllm.compose.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.localllm.compose</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/sh</string>
    <string>-c</string>
    <string>sleep 30 &amp;&amp; /opt/homebrew/bin/docker compose -f /Users/local-llm/docker/local-llm/docker-compose.yml up -d</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/var/log/llamacpp/compose.log</string>
  <key>StandardErrorPath</key>
  <string>/var/log/llamacpp/compose.error.log</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    <key>HOME</key>
    <string>/Users/local-llm</string>
    <key>DOCKER_HOST</key>
    <string>unix:///Users/local-llm/.colima/default/docker.sock</string>
  </dict>
</dict>
</plist>
