# launchagents/com.llamacpp.server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>

  <!-- Unique label for this daemon -->
  <key>Label</key>
  <string>com.llamacpp.server</string>

  <!-- Full command + arguments -->
  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/llama-server</string>
    <string>--model</string>
    <string>/opt/models/Meta-Llama-3.1-8B-Instruct-Q4_K_M.gguf</string>
    <string>--host</string>
    <string>0.0.0.0</string>
    <string>--port</string>
    <string>8080</string>
    <!-- Offload all layers to Apple Silicon GPU via Metal -->
    <string>--n-gpu-layers</string>
    <string>99</string>
    <!-- Context window size - 8192 is a good default for Home Assistant + general use -->
    <string>--ctx-size</string>
    <string>32768</string>
    <!-- Number of parallel request slots - 2 is safe for 32GB RAM with an 8B model -->
    <string>--parallel</string>
    <string>2</string>
    <!-- Keep model loaded in memory continuously -->
    <string>--keep</string>
    <string>-1</string>
  </array>

  <!-- Start at boot -->
  <key>RunAtLoad</key>
  <true/>

  <!-- Restart automatically if it crashes -->
  <key>KeepAlive</key>
  <true/>

  <!-- Run as your user so Homebrew paths and /opt/models are accessible -->
  <!-- Replace YOUR_USERNAME with the output of: whoami -->
  <key>UserName</key>
  <string>local-llm</string>

  <!-- Log stdout and stderr for debugging -->
  <key>StandardOutPath</key>
  <string>/var/log/llamacpp/server.log</string>
  <key>StandardErrorPath</key>
  <string>/var/log/llamacpp/server.error.log</string>

  <!-- Set PATH so launchd can find Homebrew binaries -->
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>

  <!-- Throttle restarts - wait 10s before restarting on failure -->
  <key>ThrottleInterval</key>
  <integer>10</integer>

</dict>
</plist>
